# infra/roles/geoip/templates/geoip-update.service.j2
[Unit]
Description=Update GeoLite2 databases
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart={{ geoip_update_script_path }}
