# infra/roles/services/go_app/templates/go-app.service.j2
# Managed by Ansible - systemd unit for {{ go_app_service_name }}
[Unit]
Description={{ go_app_service_name }} service
After={{ go_app_unit_after | join(' ') }}
Wants={{ go_app_unit_wants | join(' ') }}

[Service]
Type=simple
EnvironmentFile={{ go_app_env_file }}
{% if go_app_migrate_command %}ExecStartPre=/usr/bin/env {{ go_app_migrate_command }}
{% endif %}ExecStart=/usr/bin/env {{ go_app_exec_start }}{% if go_app_exec_args | length %} {{ go_app_exec_args | join(' ') }}{% endif %}
Restart=on-failure
RestartSec=5s
MemoryMax={{ go_app_memory_limit }}
CPUQuota={{ go_app_cpu_quota }}
WorkingDirectory={{ go_app_release_dir }}
{% if go_app_limit_nofile is not none %}
LimitNOFILE={{ go_app_limit_nofile }}
{% endif %}
{% for key, value in go_app_systemd_extra.items() %}
{{ key }}={{ value }}
{% endfor %}

[Install]
WantedBy=multi-user.target
